# package/lib/systemd/system/ahriman.timer
[Unit]
Description=ArcHlinux ReposItory MANager timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target