# package/lib/systemd/system/ahriman.service
[Unit]
Description=ArcHlinux ReposItory MANager

[Service]
ExecStart=/usr/bin/ahriman update
User=ahriman
Group=ahriman